# root/etc/s6-overlay/s6-rc.d/init-docs-sign-config/up
/etc/s6-overlay/s6-rc.d/init-docs-sign-config/run
